# trunk/library/shell/crontab
crontab -l for display
crontab -e for edit

*/1 * * * * sleep 5; /root/workspace/AnonCall/shell/loadParkedCalls5.php > /root/res.log
*/1 * * * * sleep 10; /root/workspace/tropo/shell/loadParkedCalls10.php > /root/res.log
*/1 * * * * sleep 15; /root/workspace/tropo/shell/loadParkedCalls15.php > /root/res.log
*/1 * * * * sleep 20; /root/workspace/tropo/shell/loadParkedCalls20.php > /root/res.log
*/1 * * * * sleep 25; /root/workspace/tropo/shell/loadParkedCalls25.php > /root/res.log
*/1 * * * * sleep 30; /root/workspace/tropo/shell/loadParkedCalls30.php > /root/res.log
*/1 * * * * sleep 35; /root/workspace/tropo/shell/loadParkedCalls35.php > /root/res.log
*/1 * * * * sleep 40; /root/workspace/tropo/shell/loadParkedCalls40.php > /root/res.log
*/1 * * * * sleep 45; /root/workspace/tropo/shell/loadParkedCalls45.php > /root/res.log
*/1 * * * * sleep 50; /root/workspace/tropo/shell/loadParkedCalls50.php > /root/res.log
*/1 * * * * sleep 55; /root/workspace/tropo/shell/loadParkedCalls55.php > /root/res.log
